# ops/feiniu-monitor-status.service
[Unit]
Description=Collect Feiniu Docker monitor status
After=docker.service

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/local/sbin/collect-feiniu-monitor-status.sh
