# ops/sync-html-repository.service
[Unit]
Description=Sync Feiniu Docker websites and data to GitHub
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
TimeoutStartSec=10min
ExecStart=/bin/sh /usr/local/sbin/sync-html-repository.sh
